# deploy/caldisplay.service
[Unit]
Description=CalDisplay family calendar
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=/home/erich/Desktop/CalDisplay
ExecStart=/home/erich/Desktop/CalDisplay/build/CalDisplay
Restart=always
RestartSec=2
Environment=QT_QUICK_BACKEND=software

[Install]
WantedBy=default.target
